# init/transmission
#!/sbin/openrc-run
# Copyright 1999-2011 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/net-p2p/transmission/files/transmission-daemon.initd.8,v 1.2 2011/12/04 10:35:07 swegener Exp $

extra_started_commands="reload"
description="Transmission is a fast, easy and free bittorrent client"
description_start="Start transmission-daemon server and web interface"
description_stop="Stop transmission-daemon server and web interface"
description_reload="Reload transmission-daemon settings"

rundir=/config
pidfile=/config/transmission.pid
config_dir=/config
download_dir=${TRANSMISSION_DOWNLOAD_DIR:-/data/completed}
logfile=/config/transmission.log
runas_user=transmission

SSD_OPTIONS=""

# compatibility for upgraders
TRANSMISSION_OPTIONS=${TRANSMISSION_OPTIONS:-$TD_OPTS}

check_config() {
	if [ ! -d "${rundir}" ]; then
		mkdir "${rundir}"
		if [ -n "${runas_user}" ]; then
			chown -R ${runas_user} "${rundir}"
		fi
	fi

	# In case no config directory option passed use default
	if ! $(echo ${TRANSMISSION_OPTIONS} | grep -q -e '\B-g' -e '\B--config-dir'); then
		TRANSMISSION_OPTIONS="${TRANSMISSION_OPTIONS} --config-dir ${config_dir}"
		# put download dir location on first run (and take it from config later)
		if [ ! -f ${config_dir}/settings.json ]; then
			TRANSMISSION_OPTIONS="${TRANSMISSION_OPTIONS} --download-dir ${download_dir}"
		fi
	fi

	if [ -n "${runas_user}" ]; then
		SSD_OPTIONS="${SSD_OPTIONS} --user ${runas_user}"
	fi
}

start() {
	check_config

	. /etc/transmission/environment-variables.sh

	einfo "Updating TRANSMISSION_BIND_ADDRESS_IPV4 to the ip of $IF_DEV : $IPADDR"
	export TRANSMISSION_BIND_ADDRESS_IPV4=$IPADDR

	dockerize -template /etc/transmission/settings.tmpl:config/settings.json /bin/true

	if [ ! -e "/dev/random" ]; then
		einfo "/dev/random not found - symlink to /dev/urandom"
		ln -s /dev/urandom /dev/random
	fi

	chown -R transmission:media /config
	chown transmission:media \
		${TRANSMISSION_DOWNLOAD_DIR} \
		${TRANSMISSION_INCOMPLETE_DIR} \
		${TRANSMISSION_WATCH_DIR}

	export RUN_AS=${runas_user}

	ebegin "Starting transmission daemon"
	start-stop-daemon --start --quiet --pidfile ${pidfile} ${SSD_OPTIONS} \
		--exec /usr/bin/transmission-daemon -- --pid-file ${pidfile} \
		--logfile ${logfile} ${TRANSMISSION_OPTIONS}
	eend $?

	if [ "$OPENVPN_PROVIDER" = "PIA" ]; then
		exec /etc/transmission/updatePort.sh &
	fi

}

stop() {
	ebegin "Stopping transmission daemon"
	start-stop-daemon --stop --quiet --retry TERM/45/QUIT/15 --pidfile ${pidfile}
	eend $?
}

reload() {
	ebegin "Reloading transmission configuration"
	start-stop-daemon --signal HUP --pidfile ${pidfile}
	eend $?
}

